# Spark/filesystem/etc/crontab
0 2 * * * cd /var/e2nw/scripts/; /var/e2nw/scripts/e2nwscore.rb >/dev/null 2>&1
